# dist/systemd/adblock-dnsmasq.timer
[Unit]
Description=Daily update dnsmasq configuration to block domains from blocklist

[Timer]
OnCalendar=*-*-* 03:30

[Install]
WantedBy=timers.target
